# canaryd/init_scripts/canaryd.service
[Unit]
Description=canaryd

[Service]
Type=simple
User=root
ExecStart=CANARYD_LOCATION
Restart=on-abort
